# etc/systemd/system/docker-compose-auto-update@.service
[Unit]
Description=%i auto update via docker-compose

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload docker-compose@%i.service
