# etc/systemd/system/zpool-scrub@.timer
[Unit]
Description=ZFS Scrub

[Timer]
OnCalendar=*-*-1,15 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
